# cron.hourly/slapcat-dump
#!/bin/sh -e

cd /home/dpavlin/ldap
./slapcat-list-users.pl >/dev/null
